# mng/crontab
# 一度手動で実行してkeyを登録する必要あり
* * * * * sh /root/get_log.sh
